# resources/services.d/php-fpm/finish
#!/command/execlineb -P

foreground { echo stopping php... }

exec s6-svscanctl -t /var/run/s6/services
