# resources/services.d/php-fpm/run
#!/command/execlineb -P

with-contenv

# create php socket path
foreground { mkdir -p "/var/run/php/" }

foreground { echo starting php... }

exec /usr/local/bin/php-fpm --nodaemonize
